# supervisor.conf
[program:bot]
command = ./main
autorestart = true
restart_when_binary_changed = true
[inet_http_server]
port=127.0.0.1:9001
