# supervisor.conf
[program:bot]
command = ./main
autorestart = true
restart_when_binary_changed = true
stdout_logfile = superviser_out.log
stderr_logfile = superviser_err.log
[inet_http_server]
port=127.0.0.1:30003
